# ctl/mail.service
[Unit]
Description=Start RaspbotPi for mail

[Service]
Type=simple
ExecStart=/usr/bin/python /home/fmarotta/raspbotpi/scripts/mail.py
