# ctl/raspbotpi.timer
[Unit]
Description=Timer for raspbotpi

[Timer]
OnBootSec=5min

[Install]
WantedBy=multi-user.target

